# etc/systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin mamutal91 --noclear %I $TERM
