# home/.crontab
*/20 * * * * bash $HOME/.config/scripts/dot/gitcron.sh
