# home/.crontab
*/5 * * * * $HOME/.config/scripts/gitcron.sh
*/10 * * * * cd /media/storage/Google\ Drive/ && grive
