# home/.crontab
*/10 * * * * $HOME/.config/scripts/gitcron.sh
