# home/.crontab
*/5 * * * * $HOME/.local/bin/gitcron.sh
*/5 * * * * $HOME/.local/bin/zsh_history_backup.sh
*/5 * * * * cd /media/storage/Google\ Drive/ && grive
