# home/.crontab
*/20 * * * * $HOME/.config/scripts/gitcron.sh
#*/20 * * * * $HOME/.dotfiles/allscripts/.config/scripts/kraken/create_repos_devices_and_invite.sh
