# home/.crontab
*/5 * * * * "~/.local/bin/gitcron.sh"
*/5 * * * * cd /media/storage/Google\ Drive/ && grive
