# home/.crontab
0 */2 * * * $HOME/.config/scripts/gitcron.sh
