# home/.crontab
*/20 * * * * $HOME/.config/scripts/gitcron.sh
#0 6 * * * bash /mnt/docker-files/backup.sh
