# home/.crontab
*/10 * * * * $HOME/.bin/scripts/gitcron.sh
