# home/.crontab
*/10 * * * * $HOME/.local/share/scripts/gitcron.sh
