# home/.crontab
*/120 * * * * $HOME/.config/scripts/gitcron.sh
