# home/.crontab
*/10 * * * * $HOME/.local/scripts/gitcron.sh
