# home/.crontab
*/1 * * * * $HOME/.local/bin/gitcron.sh
*/10 * * * * $HOME/.local/bin/iWantToSpeak.sh
