# exchange/crontab
*/3 * * * * $HOME/.config/scripts/zsh-history.sh
