# exchange/crontab
*/3 * * * * $HOME/.config/scripts/zsh-history.sh
*/3 * * * * $HOME/.config/scripts/gitcron.sh
*/10 * * * * $HOME/.config/scripts/iWantToSpeak.sh
